# install/systemd/kamailio.service.d/override.conf
[Unit]
After=network-online.target rtpengine-daemon.service
Wants=network-online.target

[Service]
Restart=always
RestartSec=3
LimitNOFILE=65535
LimitNPROC=65535
LimitCORE=infinity
